# runit/sv/syncthing/run
#!/bin/sh

exec syncthing -no-browser -no-restart -logfile="/home/kkga/tmp/syncthing.log" 2>&1
