# runit/sv/syncthing/run
#!/usr/bin/env bash

exec syncthing --no-browser > /dev/null
