# .local/sv/wlsunset/run
#!/bin/sh
exec wlsunset -l 55 -L 37 -t 3200 -T 6500
