# .local/sv/syncthing/run
#!/bin/sh
exec syncthing -no-browser -no-restart -logflags=0
